# reference/etc/udev/rules.d/51-ps3-ps4-disable-motion-controls.rules
#disable DS3/DS4 motion controls
SUBSYSTEM=="input", ATTRS{name}=="*Motion Sensors", RUN+="/bin/rm %E{DEVNAME}", ENV{ID_INPUT_JOYSTICK}=""