# reference/etc/udev/rules.d/99-sixpair.rules
SUBSYSTEMS=="usb", RUN+="/usr/local/bin/sixpair"